# deploy/voice-gate.service
[Unit]
Description=Voice Gate Backend (Axum)
After=network.target postgresql.service

[Service]
Type=simple
User=voice-gate
Group=voice-gate
WorkingDirectory=/opt/voice-gate
ExecStart=/opt/voice-gate/backend
Restart=always
RestartSec=5

# Environment
Environment=RUST_LOG=backend=info
Environment=HOST=0.0.0.0
Environment=PORT=3000
EnvironmentFile=/opt/voice-gate/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/voice-gate

[Install]
WantedBy=multi-user.target
